# infra/com.moskv.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  MOSKV-1 Daemon — macOS LaunchAgent
  ───────────────────────────────────
  Install:   cp launchd/com.moskv.daemon.plist ~/Library/LaunchAgents/
             launchctl load ~/Library/LaunchAgents/com.moskv.daemon.plist
  Uninstall: launchctl unload ~/Library/LaunchAgents/com.moskv.daemon.plist

  NOTE: Update the Python path below to match your environment.
        Run `which python` inside your cortex venv to get the correct path.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moskv.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <!-- ⚠️  Adjust this path to your venv Python -->
        <string>./.venv/bin/python</string>
        <string>-m</string>
        <string>cortex.daemon_cli</string>
        <string>check</string>
    </array>

    <key>WorkingDirectory</key>
    <string>.</string>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Auto-restart if daemon crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Logs: persistent location (~/.cortex/logs/) -->
    <key>StandardOutPath</key>
    <string>~/.cortex/logs/moskv-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>~/.cortex/logs/moskv-daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:./.venv/bin</string>
    </dict>
</dict>
</plist>
