# scripts/.archive/com.moskv.ouroboros-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moskv.ouroboros-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist/.venv/bin/python3</string>
        <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist/scripts/ouroboros_daily.py</string>
        <string>--verbose</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Run daily at 04:00 UTC (06:00 CEST) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist/.scratch/ouroboros/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/borjafernandezangulo/10_PROJECTS/cortex-persist/.scratch/ouroboros/launchd_stderr.log</string>

    <!-- Don't restart on failure -->
    <key>KeepAlive</key>
    <false/>

    <!-- Catch up if machine was asleep -->
    <key>StartInterval</key>
    <integer>86400</integer>
</dict>
</plist>
